# pulseAudio/conf/udev/10-udev-pulseaudio-default-source-volume.rules
# the udev rules: 
# run the script when detect the usb headset insert the system

ACTION=="add", SUBSYSTEM=="sound", RUN+="/usr/local/bin/test.sh"
